# control-api.service
[Unit]
Description=Control API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai_agents
ExecStart=/usr/bin/python3 control_api.py
Restart=always
RestartSec=10
Environment=CONTROL_API_PORT=5001
Environment=CONTROL_API_KEY=secret

[Install]
WantedBy=multi-user.target